# scripts/syncsmith.service
[Unit]
Description=Syncsmith configuration sync

[Service]
Type=oneshot
ExecStart=/opt/syncsmith/syncsmith.sh --auto
